# systemd/ollama-document-assistant.service
[Unit]
Description=Ollama Document Assistant (dry-run scanner + review web)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/ollama-document-assistant
ExecStart=%h/ollama-document-assistant/.venv/bin/python %h/ollama-document-assistant/doc_assistant_service.py --config-file %h/ollama-document-assistant/assistant_config.json
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
